# roles/gitbucket/files/gitbucket.service
[Unit]
Description="GitBucket"

[Service]
ExecStart=/bin/java -jar /gitbucket.war
Restart=always

[Install]
WantedBy=multi-user.target